# dev-util/imagefactory/files/imagefactoryd.initd
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

command="/usr/bin/imagefactoryd"
command_args="${IMAGEFACTORYD_ARGS}"
description="System Image Factory Daemon"

depend() {
    need localmount
}
